# contrib/openvpn-agent.service
[Unit]
Description=OpenVPN Management Agent
After=network.target openvpn-server@server.service
Requires=openvpn-server@server.service

[Service]
Type=simple
ExecStart=/usr/local/bin/openvpn-agent
Restart=on-failure
User=openvpn-agent
Group=openvpn-agent
ExecStart=/usr/local/bin/openvpn-agent
RuntimeDirectory=openvpn-agent
RuntimeDirectoryMode=0750
ReadWritePaths=/var/lib/openvpn/
ReadWritePaths=/run/openvpn-agent/

[Install]
WantedBy=multi-user.target
